# supervisord/rsyncd.ini
[program:rsyncd]
command=/usr/bin/rsync --daemon --no-detach
autostart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
