# supervisord/tftpd.ini
[program:tftpd]
command=/usr/sbin/in.tftpd --foreground --listen -vvv --user root --permissive --address 0.0.0.0:69 --secure --blocksize 1380 /var/lib/tftpboot
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
